# ansible/roles/relay/templates/etc/systemd/system/docker.service.d/10-args.conf.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd -H fd:// --iptables=false --bip={{ docker_net_host }} --fixed-cidr={{ docker_net_host }}
