# ansible/roles/media/files/appdata/backup-scheduler/crontab
30 3 * * * bash -c /config/backup-appdata.sh 2>&1
