# ansible/roles/common.docker/templates/etc/systemd/system/docker.service.d/10-args.conf.j2
[Service]
ExecStart=
ExecStart=/usr/bin/dockerd \
    -H fd:// \
    --bip {{ docker_net_host }} \
    --fixed-cidr {{ docker_net_host }} \
    --log-driver json-file \
    --log-opt max-file=1 \
    --log-opt max-size=200k
